# debian/cgn/etc/cron.daily/AT-bouquins2kwiki
#!/bin/sh
#
. /etc/hosts.nib.sh
if [ ! $SERVER ]; then exit ; fi
if [ `hostname` != $SERVER ]; then exit ; fi 

/var/www/bouquins/extra/bouquins2kwiki.pl --dbuser bq --dbpass 20BQ04 --dbname clt05
