# debian/cgn/etc/cron.daily/cgn_mrclean
#!/bin/sh

mrclean=/usr/bin/mrclean
LIST="/etc /root"

for dir in $LIST; do
    $mrclean $dir >> /dev/null 
done
