# debian/cgn/etc/cron.daily/cgnserv_mrtg-ink
#!/bin/sh
. /etc/hosts.nib.sh
if [ ! $SERVER ]; then exit ; fi
if [ `hostname` != $SERVER ]; then exit ; fi 

/usr/local/bin/mrtg-ink-real


