# debian/cgn/etc/cron.hourly/cgnserv_photosupdate.sh
#!/bin/sh

/usr/local/bin/cgn-photosupdate.pl


# End
