# debian/cgn/etc/cron.d/AT-apt-archives-update
09 * * * *	root	/usr/local/bin/AT-server.apt-archives-update.pl
