# debian/cgn/etc/cron.d/AT-update-bookmarks
*/5 * * * *	root	/usr/local/bin/AT-clients.update-bookmarks.pl

8 * * * *	root	/usr/local/bin/AT-server.update-bookmarks-www.pl
