# debian/cgn/etc/cron.d/AT-update-wiki
*/10 * * * *	root	/usr/local/bin/AT-server.update-wiki.pl

