# debian/cgn/etc/cron.d/cgn_updateaddress
0-50/25 * * * *	root	/usr/local/bin/cgn-updateaddress.pl
