# debian/cgn/etc/cron.d/cgn_fetchmailupdate
09 * * * *	root	/usr/local/bin/cgn-fetchmailupdate.pl
