# etc/profile.d/99.sh
# For the record:
ETC_PROFILE_SOURCED=1
[ ! -z "$DEBUG" ] && echo "/etc/profile.d sourced"
# EOF
