# etc/profile.d/stalag13-tar.sh
# tar.[bz2|gz] decompression in 4 chars. 
alias targ='tar zxvf'
alias tarb='tar jxvf'