# debian/cgn/etc/cron.weekly/AT-exim-build-stats
#!/bin/sh

if [ -e /var/log/exim4/mainlog ] && [ -e /var/www/stats ]; then
    eximstats -nr -html -charts -chartdir /var/www/stats /var/log/exim4/*  > /var/www/stats/mail.html
fi